# template/config/crontab
MAILTO=""
# m h  dom mon dow   command
* * * * * cd /var/www/html; /var/www/html/console shared:scheduler 1>/dev/null
